# 04/10-usb-keyboard.rules
ACTION=="add", ATTRS{product}=="USB KeyBoard", RUN+="/sbin/insmod /tmp/hello-1.ko"
